# motdgen.path
[Unit]
Description=Watch for /etc/coreos/motd.d changes

[Path]
PathChanged=/etc/coreos/motd.d
