# usr/lib/udev/rules.d/90-console-login-helper-messages-issuegen.rules
ACTION=="add", SUBSYSTEM=="net", ENV{INTERFACE}=="e*", RUN+="/usr/lib/console-login-helper-messages/issuegen add $env{INTERFACE}"
ACTION=="remove", SUBSYSTEM=="net", ENV{INTERFACE}=="e*", RUN+="/usr/lib/console-login-helper-messages/issuegen remove $env{INTERFACE}"
