# usr/lib/systemd/system/console-login-helper-messages-issuegen.service
[Unit]
Description=Generate console-login-helper-messages issue snippet
Documentation=https://github.com/coreos/console-login-helper-messages
Before=systemd-user-sessions.service

[Service]
Type=oneshot
# This service expects to be started only on startup.
ExecStart=/usr/libexec/console-login-helper-messages/issuegen

[Install]
WantedBy=multi-user.target
