# usr/lib/systemd/system/console-login-helper-messages-motdgen.path
[Unit]
Description=Watch for files dropped into /{etc,run,usr/lib}/console-login-helper-messages/motd.d

[Path]
PathExistsGlob=/etc/console-login-helper-messages/motd.d/*
PathExistsGlob=/run/console-login-helper-messages/motd.d/*
PathExistsGlob=/usr/lib/console-login-helper-messages/motd.d/*

[Install]
WantedBy=multi-user.target
