# usr/lib/systemd/system/issuegen.path
[Unit]
Description=Watch for /etc/console-login-helper-messages/issue.d changes

[Path]
PathChanged=/etc/console-login-helper-messages/issue.d

[Install]
WantedBy=multi-user.target
