# usr/lib/systemd/system/issuegen.path
[Unit]
Description=Watch for /run/fedora-user-messages/issue.d changes

[Path]
PathChanged=/run/fedora-user-messages/issue.d
