# usr/lib/systemd/system/issuegen.path
[Unit]
Description=Watch for /run/coreos/issue.d changes

[Path]
PathChanged=/run/coreos/issue.d
