# usr/lib/systemd/system/issuegen.path
[Unit]
Description=Watch for /etc/coreos/issue.d changes

[Path]
PathChanged=/etc/coreos/issue.d
