# usr/lib/systemd/system/motdgen.path
[Unit]
Description=Watch for /etc/fedora-user-messages/motd.d changes

[Path]
PathChanged=/etc/fedora-user-messages/motd.d
