# usr/lib/systemd/system/motdgen.path
[Unit]
Description=Watch for /run/fedora-user-messages/motd.d changes

[Path]
PathChanged=/run/fedora-user-messages/motd.d
