# usr/lib/systemd/system/motdgen.path
[Unit]
Description=Watch for /run/coreos/motd.d changes

[Path]
PathChanged=/run/coreos/motd.d
